# supervisord.conf
[supervisord]
nodaemon=true

[program:update_forecasts]
directory=/opt
command=/usr/local/bin/python2 -m pm25.rfr_forecast

[program:memcached]
command=/usr/bin/memcached -v -p 11211 -u nobody